# templates/etc/systemd/system/docker.service.d/hosts.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd://{% if docker_listen_ip != "127.0.0.1" %} -H tcp://{{ docker_listen_ip }}:{{ docker_listen_port }}{% endif %}
